# installer/updater/synapse-updater.service
[Unit]
Description=Synapse self-update daemon
Documentation=https://github.com/Iann29/convex-synapse
After=network.target docker.service
Wants=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/synapse-updater
Restart=on-failure
RestartSec=5

# Updater needs full root because setup.sh --upgrade writes to /opt/synapse,
# uses /var/run/docker.sock, and may invoke the package manager. Locking it
# down to a dedicated user would just bloat the install with sudoers config
# while not actually narrowing privilege (setup.sh ITSELF is root-equivalent).

# Defense-in-depth on the daemon itself. Notably absent: ProtectHome=true.
# It looks like a slam-dunk hardening — the daemon never reads /root or
# /home — but the daemon FORKS setup.sh which forks `docker compose up
# --build`, which calls `docker buildx` which insists on creating
# /root/.docker/ for its own config. With ProtectHome=true the unit's
# /root is read-only and buildx fails with "mkdir /root/.docker:
# read-only file system", killing every dashboard-driven upgrade.
# Operators ssh-ing in and running setup.sh manually don't see this
# because they're outside the systemd sandbox; only daemon-launched
# upgrades hit it. Real-VPS smoke caught the bug post-v1.1.0.
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

# RuntimeDirectory creates /run/synapse on each start (cleaned on stop).
# The daemon binds its unix socket inside that path; the synapse-api
# container bind-mounts the same path so it can talk to the socket
# without any TCP exposure.
RuntimeDirectory=synapse
RuntimeDirectoryMode=0755

Environment=SYNAPSE_UPDATER_SOCKET=/run/synapse/updater.sock
Environment=SYNAPSE_INSTALL_DIR=/opt/synapse

[Install]
WantedBy=multi-user.target
