# installer/updater/synapse-updater.service
[Unit]
Description=Synapse self-update daemon
Documentation=https://github.com/Iann29/convex-synapse
After=network.target docker.service
Wants=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/synapse-updater
Restart=on-failure
RestartSec=5

# Updater needs full root because setup.sh --upgrade writes to /opt/synapse,
# uses /var/run/docker.sock, and may invoke the package manager. Locking it
# down to a dedicated user would just bloat the install with sudoers config
# while not actually narrowing privilege (setup.sh ITSELF is root-equivalent).

# Defense-in-depth on the daemon itself. Notably absent: ProtectHome=true.
# It looks like a slam-dunk hardening — the daemon never reads /root or
# /home — but the daemon FORKS setup.sh which forks `docker compose up
# --build`, which calls `docker buildx` which insists on creating
# /root/.docker/ for its own config. With ProtectHome=true the unit's
# /root is read-only and buildx fails with "mkdir /root/.docker:
# read-only file system", killing every dashboard-driven upgrade.
# Operators ssh-ing in and running setup.sh manually don't see this
# because they're outside the systemd sandbox; only daemon-launched
# upgrades hit it. Real-VPS smoke caught the bug post-v1.1.0.
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

# EnvironmentFile pulls SYNAPSE_UPDATER_TOKEN (and any other shared secrets)
# from the operator's .env. Inline `Environment=` would publish the token in
# a 0644 systemd unit file readable by every local user; .env is 0600 and
# already the canonical home for this kind of credential. The leading `-`
# means "tolerate missing file" so a fresh install where the unit lands
# before .env exists fails fast at the Python layer (no token → daemon
# refuses to start) instead of failing systemd-side with a confusing
# "EnvironmentFile not found" error.
#
# {{INSTALL_DIR}} is rendered by installer/install/updater.sh::phase_install_updater
# from the operator's --install-dir flag (default /opt/synapse). Hard-coding
# /opt/synapse here (the v1.5.x bug) broke any custom install path: the
# daemon ran without SYNAPSE_UPDATER_TOKEN and refused to start, and the
# dashboard's upgrade button silently 502'd. Both EnvironmentFile and
# SYNAPSE_INSTALL_DIR resolve from the same operator flag.
EnvironmentFile=-{{INSTALL_DIR}}/.env
Environment=SYNAPSE_INSTALL_DIR={{INSTALL_DIR}}
# BIND_ADDR=0.0.0.0 (not 127.0.0.1) is REQUIRED for the synapse-api
# container to reach the daemon. Inside the container,
# `host.docker.internal` resolves to the host's docker-bridge IP
# (typically 172.17.0.1), NOT to loopback. A daemon bound to 127.0.0.1
# rejects bridge-IP connections, so /v1/admin/version_check + the
# upgrade button silently 502 from the dashboard.
#
# Security model:
#   1. Bearer token gate (every endpoint requires Authorization: Bearer
#      <SYNAPSE_UPDATER_TOKEN>). The token lives in /opt/synapse/.env
#      (0600, root-only) and is generated once at install time.
#   2. Operator firewall: port 8089 MUST be closed to the public
#      internet at the cloud-provider security-group level OR via
#      iptables/ufw. The token is defense-in-depth, not the only line.
#      Hetzner / DigitalOcean / AWS default firewalls already do this;
#      bare-metal installs without ufw need explicit `ufw deny 8089`.
Environment=SYNAPSE_UPDATER_BIND=0.0.0.0

[Install]
WantedBy=multi-user.target
